# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/backup-etcd.service
[Unit]
Description=Backup ETCD
[Service]
Type=oneshot
ExecStart=/opt/app/current/bin/node/etcd-backup.sh