# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/backup-etcd.timer
[Unit]
Description=Timer to backup ETCD
[Timer]
OnUnitActiveSec=6h
Unit=backup-etcd.service
[Install]
WantedBy=multi-user.target