# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/hostnic-status.timer
[Unit]
Description=update hostnic status every minute
ConditionFileNotEmpty=/data/appctl/data/hostnic.status

[Timer]
OnActiveSec=1
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=multi-user.target
