# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/hostnic-status.service
[Unit]
Description=update hostnic status

[Service]
Type=simple
ExecStart=/usr/bin/appctl updateHostnicStatus
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
