# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/ks-installer.service
[Unit]
Description=install KubeSphere

[Service]
ExecStart=/usr/bin/appctl setUpKs
ExecStopPost=/usr/bin/appctl reloadKsEip

[Install]
WantedBy=multi-user.target
