# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/kube-lb.service
[Unit]
Description=create kube-apiserver load balancer

[Service]
ExecStart=/usr/bin/appctl setUpKubeLb
ExecStopPost=/usr/bin/appctl distributeKubeLbFile

[Install]
WantedBy=multi-user.target
