# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/vbc-nic.service
[Unit]
Description=set up VBC nic

[Service]
Type=simple
ExecStart=/usr/bin/appctl setUpVbcNic
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
