# ansible/roles/app-role-k8s/files/opt/app/current/conf/systemd/docker-prune.service.d/10-appctl.conf
[Service]
ExecStart=
ExecStart=/usr/bin/appctl runDockerPrune
