# ansible/roles/k8s-node/files/opt/app/current/conf/systemd/docker-prune.service
[Unit]
Description=run docker system prune

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -f

[Install]
WantedBy=multi-user.target
