# ansible/roles/k8s-node/files/opt/app/current/conf/systemd/kube-certs.service
[Unit]
Description=check and rotate kube certs if neccessary

[Service]
Type=oneshot
ExecStart=/usr/bin/kubeadm alpha certs renew all

[Install]
WantedBy=multi-user.target
